# build_scripts/tools/rc.local
#!/bin/sh -e

/etc/getCert.sh &

echo 1 > /root/lava

exit 0
